# travel-backend.service
[Unit]
Description=Travel Management Backend API
Documentation=https://github.com/yourusername/travel
After=network.target postgresql.service

[Service]
Type=simple
# USER YANG VALID: Ganti dengan output dari: whoami
# Atau comment baris User untuk run as root (testing only)
User=tazkia
WorkingDirectory=/home/tazkia/travel/backend
Environment=NODE_ENV=production
Environment=PORT=5000
Environment=DATABASE_URL=postgresql://user:password@localhost:5432/travel_db
Environment=JWT_SECRET=your-jwt-secret-here

# Path to node and your app
ExecStart=/usr/bin/node src/server.js

# Restart policy
Restart=always
RestartSec=10

# Performance
LimitNOFILE=65536

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=travel-backend

[Install]
WantedBy=multi-user.target
